# cex_microstructure_mvd/scripts/systemd/cex-diskmon.timer
[Unit]
Description=Run CEX disk-usage monitor hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
